# mining-main/build/mining-deployment/amazon-web-services/internal/Resources/MiningApi/Services/miningApiStart.service
[Unit]
Description=mining-api server start service
After=network.target
[Service]
Restart=always
User=ec2-user
ExecStart=/home/ec2-user/server/start.sh

[Install]
WantedBy=default.target